# .crontab
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/home/death/dev/bin

# volume, console keyboard repeat rate
@reboot sudo aumix -v 80 -c 0 -i 0 -l 0 -p 0 -w 80 -1 0 -2 0 -3 0 >/dev/null 2>&1
@reboot sudo kbdrate -d 210 -r 83

# kernel tweaks
@reboot sudo su -c 'echo 25 > /proc/sys/vm/swappiness'
@reboot sudo su -c 'echo 1500 > /proc/sys/vm/dirty_writeback_centisecs'

# clean up crap infrequenty
0 0 * * * DoItEasy rm -rf $HOME/.wine/drive_c/windows/temp/*
0 0 * * * DoItEasy find $HOME/.vim/backup -type f -mtime +30 | xargs rm -f
0 0 * * * DoItEasy sudo apt-get clean

# and then fetch gmail
10 0 * * * while DoItEasy fetchmail -kUa -m maildrop; do :; done;

# and then
0 1 * * * Backup
0 4 * * sat Offsite

#@hourly export DISPLAY=:0; test `xprintidle` -gt $((1000*60*60*4)) && python $HOME/dev/mpd-weighted-playlist/mpd_weighted_playlist.py < $HOME/.mpd-weights >/dev/null

# clean up small crap frequently
@hourly rmdir $HOME/Desktop 2>/dev/null

# clean up (maybe) big crap frequently
@hourly DoItEasy rm -f $HOME/core $HOME/core.[0123456789]* 2>/dev/null
@hourly DoItEasy rm -rf $HOME/.local/share/Trash 2>/dev/null

# dump this crontab for tracking in git
@daily crontab -l > $HOME/.crontab

